# archlinux/eyed.service
[Unit]
Description=pam_eye server gateway

[Service]
ExecStart=/usr/bin/eyed
Restart=always

[Install]
WantedBy=multi-user.target
